# server/config/supervisor_template.conf
[program:PROJECT_NAME]
command = REPO_DIR/scripts/exec_gunicorn.sh
user = SERVER_USER
stdout_logfile = REPO_DIR/run/gunicorn.log
redirect_stderr = true
autorestart = true
